# root/etc/periodic/weekly/mastodon_cleanup
#!/bin/sh

RAILS_ENV=production 

/app/www/bin/tootctl media remove
/app/www/bin/tootctl preview_cards remove